# deploy/systemd/alertbox-scheduler.service
[Unit]
Description=AlertBox Scheduler
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/alertbox
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/alertbox/.venv/bin/python -m app.scheduler
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
